# tools/cursor-bridge/com.cursor-bridge.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cursor-bridge</string>

  <key>ProgramArguments</key>
  <array>
    <string>__VENV_PYTHON__</string>
    <string>__BRIDGE_PY__</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>StandardOutPath</key>
  <string>__LOG_DIR__/cursor-bridge.out.log</string>

  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/cursor-bridge.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>CURSOR_BRIDGE_SOCKET</key>
    <string>__SOCKET_PATH__</string>
    <key>CURSOR_BRIDGE_LOG</key>
    <string>__LOG_DIR__/cursor-bridge.log</string>
  </dict>
</dict>
</plist>
